# dotagents/tmux-agent-status/com.agents.agent-status.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agents.agent-status</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/agent-status</string>
        <string>daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/agent-status.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/agent-status.log</string>
</dict>
</plist>
