# configs/cron/codex-prefix-thread-projects.cron
PATH=/home/mikeastock/.local/bin:/home/mikeastock/.local/share/mise/installs/node/latest/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin

* * * * * codex-prefix-thread-projects --detect-merged >>/tmp/codex-prefix-thread-projects.log 2>&1
